# deploy/systemd/agentrt-hook.service
# AgentRT systemd unit 文件 — Layer 2: Agent 服务
# Copyright (C) 2025-2026 SPHARX Ltd.
# P1.23.2: daemon 启动顺序编排

[Unit]
Description=AgentRT Hook Daemon
After=network.target agentrt-tool.service
Requires=agentrt-tool.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/hook_d
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=@AGENTOS_RUNTIME_DIR@
LimitNPROC=512

[Install]
WantedBy=multi-user.target
